# lib/systemd/system/failure-email@.service
[Unit]
Description=OnFailure for %i

[Service]
Type=oneshot
ExecStart=/usr/sbin/systemd_onfailure_email.sh %i